# systemd/user/vocoder.service
[Unit]
Description=Vocoder dictation daemon (Wayland-safe, local-only)
Documentation=local
After=graphical-session.target

[Service]
Type=simple
# --- Environment (override via your shell/env if desired)
Environment=PYTHONUNBUFFERED=1
Environment=WHISPER_API_URL=http://127.0.0.1:8767/v1/transcribe
Environment=WHISPER_LANG=en
Environment=WHISPER_DIARIZE=false
Environment=WHISPER_NUM_SPEAKERS=2
Environment=WHISPER_FORMAT=json
# Optional: set a custom socket path
# Environment=VOCODER_SOCK=%h/.cache/vocoder.sock

# --- Daemon entrypoint (uses your per-user venv)
ExecStart=%h/whisper-dictation/.venv/bin/python %h/whisper-dictation/vocoderd.py
Restart=on-failure
RestartSec=0.5

# Hardening (still user-scoped; keep simple for audio access)
NoNewPrivileges=true

[Install]
WantedBy=default.target